# ansible/files/fileserver.service
[Unit]
Description=App

[Service]
ExecStart=/home/ubuntu/bin/main_fileserver

[Install]
WantedBy=multi-user.target